# DjangoBlog.conf
[program:DjangoBlog]
command = gunicorn DjangoBlog.wsgi_prod -b 127.0.0.1:8081
directory = /root/DjangoBlog/DjangoBlog
autostart = true
autorestart = true
